# conf/supervisord.conf
[supervisord]
nodaemon=true

[program:httpd]
command=/usr/sbin/httpd
autostart=true
autorestart=true

[program:php-fpm]
command=/usr/sbin/php-fpm
autostart=true
autorestart=true
